# roles/flannel/templates/cni-bridge.service
[Unit]
Requires=network.target
Before=docker.service
[Service]
Type=oneshot
ExecStart=/usr/local/bin/cni-bridge.sh
RemainAfterExit=true